# buildroot-kiosk/board/sdlw/raspberrypi/rootfs_overlay/etc/profile.d/xdg-runtime-dir
export XDG_RUNTIME_DIR=/tmp/.runtime-$USER

mkdir -p $XDG_RUNTIME_DIR
chmod 700 $XDG_RUNTIME_DIR
